# deploy/supervisor/hyperlit-embeddings.conf
; ─────────────────────────────────────────────────────────────────────────────
; Dedicated worker for the `embeddings` queue (GenerateNodeEmbedding — one short
; HTTP call per node; QueueBookEmbeddings fans a book out into thousands of them).
;
; WHY ITS OWN WORKER: embeddings previously rode along on hyperlit-worker's
; `--queue=default,embeddings` list. Priority order protected imports from
; *queued* embeddings, but a worker BUSY on an embedding job is still a worker
; an import has to wait for — and during a backlog drain, import capacity
; effectively halves. Per-class isolation (same as citation/vibe) ends that.
; The worker itself is light (~50–80 MB, measured tests/load/memprobe.sh).
;
; ⚠️ Ship together with the hyperlit-worker.conf whose --queue list is `default`
; only — if the old conf still says `default,embeddings`, both serve this queue.
;
; Install on the droplet (app dir /var/www/hyperlit, SSH marx@170.64.145.89):
;   sudo cp deploy/supervisor/hyperlit-embeddings.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update
;
; --max-jobs=1000: embedding jobs are tiny; recycle the PHP process often enough
; that a long backlog drain can't accumulate leaked memory.
; ─────────────────────────────────────────────────────────────────────────────
[program:hyperlit-embeddings]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/hyperlit/artisan queue:work --queue=embeddings --sleep=3 --tries=2 --timeout=300 --max-jobs=1000
directory=/var/www/hyperlit
user=marx
numprocs=1
autostart=true
autorestart=true
startsecs=1
stopwaitsecs=310
stopsignal=SIGTERM
redirect_stderr=true
stdout_logfile=/var/www/hyperlit/storage/logs/embeddings-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
