# deploy/supervisor/hyperlit-citation.conf
; ─────────────────────────────────────────────────────────────────────────────
; Dedicated worker for the `citation-pipeline` queue (CitationPipelineJob,
; CitationScanBibliographyJob, CanonicalizeLibraryJob).
;
; WHY ITS OWN WORKER: a citation pipeline run is 12–15+ minutes of LLM calls and
; web fetches (job $timeout = 7200s). It used to share the serial hyperlit-worker
; with `default` — and OUTRANKED it, so one user's citation review head-of-line-
; blocked every other user's document import for the whole run. Same disease,
; same cure as hyperlit-vibe.conf: own queue, own worker, real parallelism.
;
; ⚠️ Ship together with the hyperlit-worker.conf that DROPS citation-pipeline
; from its --queue list. If you install this without updating hyperlit-worker,
; both serve the queue (harmless but defeats the isolation); if you update
; hyperlit-worker without installing this, citation jobs queue forever.
;
; ⚠️ Also requires the config/queue.php retry_after bump (7500s). At the old 90s
; a >90s job gets re-reserved by any parallel worker and runs twice.
;
; Install on the droplet (app dir /var/www/hyperlit, SSH marx@170.64.145.89):
;   sudo cp deploy/supervisor/hyperlit-citation.conf /etc/supervisor/conf.d/
;   sudo cp deploy/supervisor/hyperlit-worker.conf  /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update
;   php artisan queue:restart
;
; TIMEOUT ORDERING (must hold so nothing kills a job mid-run):
;   job $timeout 7200s  ≤  worker --timeout 7200s  <  stopwaitsecs 7210s
;   stopwaitsecs > timeout means `supervisorctl restart` waits for an in-flight
;   run instead of SIGKILLing it — prefer `php artisan queue:restart` on deploys
;   so you never actually wait the full 2 h.
;
; SourceNetworkHarvestJob honors the 7200s cap PER SLICE, not per harvest: the
; runner stops at its slice budget (source_harvest.slice_seconds, default 900s)
; and the job re-dispatches itself, so an arbitrarily long harvest is a chain
; of short jobs. Deploy-friendly: `php artisan queue:restart` waits at most one
; slice (+ the in-flight work), never hours. A crashed or timed-out slice still
; finalizes the yield report from persisted state (failed() →
; HarvestRunner::finalize, marked partial), and a user re-trigger resumes from
; the persisted bookmark.
;
; RAM (~1.9 GB + 2 GB swap, OOM history): citation work is HTTP/LLM-bound, not
; memory-hungry like imports (no pandoc/Python OCR), so one extra PHP worker
; (~60–80 MB) is safe where numprocs=2 on imports is not. Keep numprocs=1; the
; lever for parallel citation runs is numprocs here, RAM permitting.
; ─────────────────────────────────────────────────────────────────────────────
[program:hyperlit-citation]
process_name=%(program_name)s_%(process_num)02d
; xvfb-run -a: the OA fetch ladder's last rung is a HEADED patchright browser
; (headless loses Cloudflare managed challenges — see deploy/oa-fetch-hardening.md).
; A screenless droplet has no display, so headed Chrome must draw into xvfb's
; fake framebuffer. Without this prefix the browser lane dies `browser_launch_failed`.
; Requires: sudo apt install -y xvfb xauth  (xvfb-run shells out to xauth).
command=xvfb-run -a php /var/www/hyperlit/artisan queue:work --queue=citation-pipeline --sleep=3 --tries=1 --timeout=7200 --max-jobs=10
directory=/var/www/hyperlit
user=marx
; HOME: supervisor does NOT set it for the worker, and headed Chromium needs it
; — its crashpad handler derives --database from $HOME and aborts on launch when
; unset (`chrome_crashpad_handler: --database is required` => the browser dies =>
; browser_launch_failed). Must match user= above.
; PLAYWRIGHT_BROWSERS_PATH: pin the patchright browser to a fixed path (NOT the
; default per-user ~/.cache — with an unset/wrong HOME that resolves to an empty
; dir => `Executable doesn't exist`).
; Install once: sudo PLAYWRIGHT_BROWSERS_PATH=/opt/ms-playwright npx patchright install chromium
environment=HOME="/home/marx",PLAYWRIGHT_BROWSERS_PATH="/opt/ms-playwright"
numprocs=1
autostart=true
autorestart=true
startsecs=1
stopwaitsecs=7210
stopsignal=SIGTERM
redirect_stderr=true
stdout_logfile=/var/www/hyperlit/storage/logs/citation-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
