# deploy/supervisor/hyperlit-audio.conf
; ─────────────────────────────────────────────────────────────────────────────
; Dedicated worker for the `audio` queue (GenerateBookAudioJob — per-node TTS
; audiobook generation, docs/audio.md).
;
; WHY ITS OWN WORKER: a full-book synthesis is hundreds→thousands of TTS API
; calls and holds the worker for minutes (job $timeout 3600s). On `default` it
; would head-of-line-block every user's import behind one book's narration —
; the exact lesson already learned twice (vibe, citation). Same pattern as
; hyperlit-vibe.
;
; Install on the droplet (app dir /var/www/hyperlit, SSH marx@170.64.145.89):
;   sudo cp deploy/supervisor/hyperlit-audio.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl start hyperlit-audio:*
;
; TIMEOUT ORDERING (must hold so nothing kills a job mid-run):
;   job $timeout 3600s  =  worker --timeout 3600s  <  stopwaitsecs 3610s
;   and retry_after (config/queue.php, 7500s) > 3600s so a long book is never
;   re-reserved and double-billed by a parallel worker.
;
; INTERRUPTION IS CHEAP: the job is resumable — it skips every node whose
; (node_id, source_hash) audio already exists and bills only what it actually
; synthesized, so a SIGKILL mid-book loses nothing but the in-flight batch;
; the user just presses generate again. stopwaitsecs still waits politely.
;
; RAM: light — a PHP worker doing batched HTTP calls (concurrency 5, ~150 KB
; base64 audio per response); embeddings-class footprint (~50 MB est).
; numprocs=1 — one book narrates at a time; the per-book Cache::lock also
; serializes duplicate requests upstream. --max-jobs recycles the PHP process
; so leaked memory can't accumulate across books.
; ─────────────────────────────────────────────────────────────────────────────
[program:hyperlit-audio]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/hyperlit/artisan queue:work --queue=audio --sleep=3 --tries=1 --timeout=3600 --max-jobs=20
directory=/var/www/hyperlit
user=marx
numprocs=1
autostart=true
autorestart=true
startsecs=1
stopwaitsecs=3610
stopsignal=SIGTERM
redirect_stderr=true
stdout_logfile=/var/www/hyperlit/storage/logs/audio-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
