# deploy/supervisor/hyperlit-vibe.conf
; ─────────────────────────────────────────────────────────────────────────────
; Dedicated worker for the `vibe` queue (VibeConversionJob).
;
; WHY ITS OWN WORKER: a vibe conversion shells out to Python for up to ~28 min
; (the job's Symfony Process timeout is 1700s). It used to sit on the `default`
; queue alongside imports (ProcessDocumentImportJob), and `default` is served by
; ONE serial worker (hyperlit-worker) — so a single vibe run head-of-line-blocked
; every user's import for up to half an hour. Moving it to its own queue + its own
; worker lets vibe and imports run in PARALLEL. (Same pattern as hyperlit-citation.)
;
; Install on the droplet (app dir /var/www/hyperlit, SSH marx@170.64.145.89):
;   sudo cp deploy/supervisor/hyperlit-vibe.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl start hyperlit-vibe:*
;
; TIMEOUT ORDERING (must hold so nothing kills a job mid-run):
;   Process timeout 1700s  <  worker --timeout 1800s  <  stopwaitsecs 1810s
;   → stopwaitsecs > timeout means `supervisorctl restart` waits for an in-flight
;     vibe run to finish gracefully instead of SIGKILLing it.
;
; DOCKER: if VIBE_SANDBOX_IMAGE is set in .env, this worker runs `docker` — the
; `marx` user must be in the `docker` group (`sudo usermod -aG docker marx`).
;
; RAM (~1.9 GB + 2 GB swap, OOM history): keep numprocs=1 — only ever ONE vibe
; run at a time. --max-jobs recycles the PHP worker periodically so leaked memory
; can't accumulate across runs. Match `user=` to your existing hyperlit-worker.conf.
; ─────────────────────────────────────────────────────────────────────────────
[program:hyperlit-vibe]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/hyperlit/artisan queue:work --queue=vibe --sleep=3 --tries=1 --timeout=1800 --max-jobs=20
directory=/var/www/hyperlit
user=marx
numprocs=1
autostart=true
autorestart=true
startsecs=1
stopwaitsecs=1810
stopsignal=SIGTERM
redirect_stderr=true
stdout_logfile=/var/www/hyperlit/storage/logs/vibe-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
