# deploy/supervisor/hyperlit-worker.conf
; ─────────────────────────────────────────────────────────────────────────────
; Main worker — the `default` queue (ProcessDocumentImportJob = doc imports &
; reconverts) plus `citation-pipeline`. This is a REFERENCE template; reconcile it
; with the hyperlit-worker.conf already on the droplet (match user=, log paths).
;
; IMPORT CONCURRENCY — the `numprocs` lever:
;   numprocs=N runs N copies of this worker, so up to N imports run AT ONCE. With
;   numprocs=1 (current) two users importing simultaneously serialize — the second
;   waits for the first. Raising numprocs is how you let them overlap.
;
;   BUT this droplet is ~1.9 GB RAM + 2 GB swap with a history of OOM-killing
;   PHP-FPM (→ Cloudflare 502). An import runs pandoc + Python (+ sometimes Mistral
;   OCR) and can spike to several hundred MB. Two concurrent imports + the vibe
;   worker + PHP-FPM can exhaust RAM and fall into swap (slow) or OOM (502s).
;   → Keep numprocs=1 until RAM is increased. The biggest win (vibe no longer
;     blocking imports) comes from hyperlit-vibe.conf, NOT from more import workers.
;     If concurrent imports become common, vertically scale the droplet first, THEN
;     raise numprocs to 2. --max-jobs recycles workers so memory can't creep.
;
; If you split citation-pipeline into its own hyperlit-citation.conf already, drop
; `citation-pipeline,` from the --queue list below to avoid double-serving it.
; ─────────────────────────────────────────────────────────────────────────────
[program:hyperlit-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/hyperlit/artisan queue:work --queue=citation-pipeline,default --sleep=3 --tries=1 --timeout=900 --max-jobs=50
directory=/var/www/hyperlit
user=marx
numprocs=1
autostart=true
autorestart=true
startsecs=1
stopwaitsecs=910
stopsignal=SIGTERM
redirect_stderr=true
stdout_logfile=/var/www/hyperlit/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
