# deploy/supervisor/hyperlit-worker.conf
; ─────────────────────────────────────────────────────────────────────────────
; Main worker — the `default` queue ONLY (ProcessDocumentImportJob = doc imports
; & reconverts). This is a REFERENCE template; reconcile it with the
; hyperlit-worker.conf already on the droplet (match user=, log paths).
;
; `citation-pipeline` moved to its OWN worker (hyperlit-citation.conf) — it used
; to live here ahead of `default` in the --queue list, so a 15-min citation
; review head-of-line-blocked every import. Document conversion must never wait
; behind any other job class: every queue gets its own Supervisor program.
; ⚠️ Ship this together with hyperlit-citation.conf (see its header) and the
; config/queue.php retry_after bump (7500s — at 90s, parallel workers re-reserve
; any >90s job and run it twice).
;
; IMPORT CONCURRENCY — the `numprocs` lever:
;   numprocs=N runs N copies of this worker, so up to N imports run AT ONCE. With
;   numprocs=1 (current) two users importing simultaneously serialize — the second
;   waits for the first. Raising numprocs is how you let them overlap ("standby
;   workers" for multi-user load).
;
;   BUT this droplet is ~1.9 GB RAM + 2 GB swap with a history of OOM-killing
;   PHP-FPM (→ Cloudflare 502). MEASURED (memprobe.sh, 2026-06-12): one import
;   peaks at 212 MB; all four job classes simultaneously ~645 MB worst case —
;   fits at numprocs=1. Each extra import worker adds another 212 MB at peak,
;   which pushes worst case to the edge of physical RAM on this box.
;   Full RAM-budget math: deploy/supervisor/README.md.
;   → Keep numprocs=1 until RAM is increased. Order of operations: vertically
;     scale the droplet first, THEN raise numprocs to 2. --max-jobs recycles
;     workers so memory can't creep.
; ─────────────────────────────────────────────────────────────────────────────
[program:hyperlit-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/hyperlit/artisan queue:work --queue=default --sleep=3 --tries=1 --timeout=900 --max-jobs=50
directory=/var/www/hyperlit
user=marx
numprocs=1
autostart=true
autorestart=true
startsecs=1
stopwaitsecs=910
stopsignal=SIGTERM
redirect_stderr=true
stdout_logfile=/var/www/hyperlit/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
