# deploy/supervisor/hyperlit-search.conf
; ─────────────────────────────────────────────────────────────────────────────
; Dedicated worker for the `search-supplement` queue
; (IngestExternalCitationCandidatesJob — OpenAlex/Open Library fetch when a
; citation search returns thin local results).
;
; WHY ITS OWN WORKER: a user is ACTIVELY WAITING on this job — the citation
; modal polls for ~8s after dispatch ("searching external databases…"). On the
; shared `default` worker it queues behind document imports (up to 15 min), so
; the polls always miss and the UX degrades to "results appear whenever".
; Same disease, same cure as hyperlit-citation / hyperlit-vibe: own queue,
; own worker.
;
; ⚠️ Ships together with the job's onQueue('search-supplement') — README
; invariant #1: a queue nobody listens on silently never runs. Install:
;   sudo cp deploy/supervisor/hyperlit-search.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update
;   php artisan queue:restart
;
; RAM (~1.9 GB box, budget in deploy/supervisor/README.md): this is the
; lightest job class we have — two HTTP fetches + ~10 small upserts, well
; under the embeddings worker's measured 50 MB peak. One extra PHP worker
; (~50–60 MB peak, ~35 MB idle) fits the current budget with the same margin
; the citation worker was approved on.
;
; --sleep=1 (not 3): pickup latency is user-visible here — the modal's first
; poll lands 2.5s after dispatch. One extra idle DB poll per second is noise.
; TIMEOUT ORDERING: job $timeout 120s ≤ worker --timeout 120s < stopwaitsecs 130s.
; ─────────────────────────────────────────────────────────────────────────────
[program:hyperlit-search]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/hyperlit/artisan queue:work --queue=search-supplement --sleep=1 --tries=1 --timeout=120 --max-jobs=200
directory=/var/www/hyperlit
user=marx
numprocs=1
autostart=true
autorestart=true
startsecs=1
stopwaitsecs=130
stopsignal=SIGTERM
redirect_stderr=true
stdout_logfile=/var/www/hyperlit/storage/logs/search-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
